# raspberry-pi/deploy/solar_tracking/對照組1/solar_tracking.service
[Unit]
Description=Solar Tracking Controller (對照組1)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/solar_tracking/對照組1
ExecStart=/usr/bin/python3 /home/pi/solar_tracking/對照組1/traditional_controller.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/pi/solar_tracking/對照組1/service.log
StandardError=append:/home/pi/solar_tracking/對照組1/service.log

[Install]
WantedBy=multi-user.target
