# raspberry-pi/deploy/solar_tracking/對照組2/solar_tracking.service
[Unit]
Description=Solar Tracking Controller (對照組2, system_id=4)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/solar_tracking/對照組2
ExecStart=/usr/bin/python3 /home/pi/solar_tracking/對照組2/traditional_controller.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/pi/solar_tracking/對照組2/service.log
StandardError=append:/home/pi/solar_tracking/對照組2/service.log

[Install]
WantedBy=multi-user.target