# raspberry-pi/deploy/solar_tracking/實驗組2/solar_tracking.service
[Unit]
Description=Solar Tracking Controller (т»джЕЌухё2, system_id=2)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/solar_tracking/т»джЕЌухё2
ExecStart=/usr/bin/python3 /home/pi/solar_tracking/т»джЕЌухё2/anfis_controller.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/pi/solar_tracking/т»джЕЌухё2/service.log
StandardError=append:/home/pi/solar_tracking/т»джЕЌухё2/service.log

[Install]
WantedBy=multi-user.target