# raspberry-pi/deploy/solar_tracking/實驗組1/solar_tracking.service
[Unit]
Description=Solar Tracking Controller (т»джЕЌухё1)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/solar_tracking/т»джЕЌухё1
ExecStart=/usr/bin/python3 /home/pi/solar_tracking/т»джЕЌухё1/anfis_controller.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/pi/solar_tracking/т»джЕЌухё1/service.log
StandardError=append:/home/pi/solar_tracking/т»джЕЌухё1/service.log

[Install]
WantedBy=multi-user.target
